# systemd/sirius-adminbot.service
[Unit]
Description=Sirius Loyalty System Admin Bot
After=network.target sirius-server.service

[Service]
Type=simple
User=goomer125
WorkingDirectory=/opt/sirius
Environment="CONFIG_PATH=/opt/sirius/production.yaml"
Environment="TOKEN_PATH=/opt/sirius/config/admin_token.txt"
Environment="ADMINS_PATH=/opt/sirius/config/admins.json"
Environment="SERVER_URL=http://localhost:8080"
Environment="API_TOKEN_PATH=/opt/sirius/config/api_token.txt"
ExecStart=/opt/sirius/bin/adminbot --config /opt/sirius/production.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sirius-adminbot

[Install]
WantedBy=multi-user.target
